# config/supervisor_kolibri.template.conf
[program:kolibri]
directory = {{KOLIBRI_HOME}}
command = {{KOLIBRI_HOME}}/startkolibri.sh
user = {{KOLIBRI_USER}}
stdout_logfile = {{KOLIBRI_HOME}}/kolibri_supervisor.log
redirect_stderr = true
